# ubuntu-noble/_root/etc/s6-overlay/s6-rc.d/userpasswd/up
/etc/s6-overlay/scripts/set-user-passwd.sh
