# ubuntu-noble/_root/etc/s6-overlay/s6-rc.d/novnc/run
#!/command/with-contenv sh

s6-setuidgid user
/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 0.0.0.0:8080
