# ubuntu-noble/_root/etc/s6-overlay/s6-rc.d/tigervnc/run
#!/command/with-contenv sh

s6-setuidgid user /usr/bin/Xvnc "$DISPLAY" -geometry "${SCR_WIDTH}x${SCR_HEIGHT}" -depth 32 -SecurityTypes none -nolock
